# deploy/heartpet.service.example
[Unit]
Description=HeartPet Tierverwaltung
After=network.target

[Service]
Type=simple
# Anpassen, z. B. /opt/HeartPet oder /root/HeartPet
WorkingDirectory=/opt/HeartPet
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=5
# Anpassen, z. B. www-data oder root
User=www-data
Group=www-data
Environment=PORT=3000
Environment=HEARTPET_SESSION_SECRET=bitte-aendern

[Install]
WantedBy=multi-user.target
